# platform/infra/systemd.service
[Unit]
Description=Monitoring System Service
After=network.target

[Service]
User=runner
ExecStart=/usr/bin/python /app/monitoring_system.py
Restart=always

[Install]
WantedBy=multi-user.target
